# pkg/pkgroot/Library/LaunchAgents/com.github.groob.yo-yo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Disabled</key>
        <false/>
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        </dict>
        <key>Label</key>
        <string>com.github.groob.yo-yo</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/yo-yo</string>
            <string>-n</string>
            <string>/Library/Management/yo-yo</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>WatchPaths</key>
        <array>
            <string>/Library/Management/yo-yo/</string>
        </array>
    </dict>
</plist>
